# layout/System/Library/LaunchDaemons/ws.hbang.aptdated.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ws.hbang.aptdated</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/aptdated</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StartInterval</key>
  <integer>7200</integer>
  <key>UserName</key>
  <string>mobile</string>
  <key>StandardErrorPath</key>
  <string>/tmp/aptdated.log</string>
  <key>StandardOutPath</key>
  <string>/tmp/aptdated.log</string>
</dict>
</plist>
